# codes/source_code/01_pc_tools/board_ch341_autobind.rules
ACTION=="add", SUBSYSTEM=="usb_device", RUN+="/etc/udev/ch341-autobind.sh"
